# deploy/prop-paper-trading.timer
[Unit]
Description=Run Prop Firm Paper Trading on Weekdays

[Timer]
OnCalendar=Mon..Fri 03:55:00 America/New_York
Persistent=true

[Install]
WantedBy=timers.target
